# scripts/podcst-poller.service
[Unit]
Description=Podcst feed poller daemon
After=network.target network-online.target postgresql.service
Requires=network-online.target

[Service]
Type=simple
EnvironmentFile=/home/shantanu/src/shantanuraj/podcst-web/.env
ExecStart=/home/shantanu/.bun/bin/bun run /home/shantanu/src/shantanuraj/podcst-web/scripts/poll-feeds.ts --daemon
WorkingDirectory=/home/shantanu/src/shantanuraj/podcst-web
Restart=always
RestartSec=10
User=shantanu
Group=shantanu
StandardOutput=journal
StandardError=journal
SyslogIdentifier=podcst-poller

[Install]
WantedBy=multi-user.target
